# root/etc/s6-overlay/s6-rc.d/svc-geoip2influx/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

  exec \
        python3 /config/geoip2db/run.py